# crontab
0 7 * * * cd /usr/src/app && python ficker.py >> log.txt 2>&1