# crontab
0 8 * * * cd /usr/src/app && python ficker.py >> log.txt 2>&1

