# crontab
* * * * * python3 /usr/src/app/ficker.py
* * * * * touch /test.txt
