# ficker.cron
0 6 * * * cd /root && python3 ficker.py >> /var/log/cron.log 2>&1
